# bbs.conf
[program:bbs]
command=/var/www/venv/bin/gunicorn wsgi --bind 0.0.0.0:2000 --pid /tmp/bbs.pid
directory=/var/www/bbs
autostart=true
autorestart=true