# bbs.conf
[program:bbs]
command=/var/www/venv/gunicorn wsgi -c gunicorn.config.py
directory=/var/www/bbs
autostart=true
autorestart=true